# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:artisan-serve]
command=php /var/www/html/artisan serve --host=0.0.0.0 --port=8000
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/artisan-serve.log
stderr_logfile=/var/log/artisan-serve.err

[program:npm-dev]
command=npm run dev
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/npm-dev.log
stderr_logfile=/var/log/npm-dev.err


[program:es-index-eav]
command=bash -c "while true; do php /var/www/html/artisan es:index-eav || true; sleep 30; done"
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/var/log/es-index.log
stderr_logfile=/var/log/es-index.err
